# supervisor-sample.conf
[program:main]
command=python main.py
directory=/home/deploy/
autostart=true
autorestart=true
stdout_logfile=/home/deploy/main.log
redirect_stderr=true
